# config/supervisord/conf.d/pilbox.conf
[program:pilbox]
command = /path/to/bin/pilbox --config=/path/to/pilbox_config.py --debug
autostart = true
autorestart = true
user = system_user_name
